# systemwide/udev-hdmi-hotplug.rules
SUBSYSTEM=="drm", ACTION=="change", RUN+="/bin/bash /usr/local/bin/autorandr-detect"
